# avbuzz.service
[Unit]
Description=AVBUZZ Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/avbuzz
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
